# systemd/pi-menu-server.service
[Unit]
Description=Pi Menu Display Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pi-menu-display
ExecStart=/home/pi/pi-menu-display/venv/bin/uvicorn server.app:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
